# systemd/birdstation-export.timer
[Unit]
Description=Birdstation — export every 30 minutes

[Timer]
OnBootSec=5min
OnUnitActiveSec=30min
Unit=birdstation-export.service

[Install]
WantedBy=timers.target
