# systemd/birdstation-export.service
[Unit]
Description=Birdstation — export parquet to R2 and trigger Vercel deploy
After=network.target

[Service]
Type=oneshot
User=birdstation
WorkingDirectory=/opt/birdstation
EnvironmentFile=/opt/birdstation/.env
ExecStart=/opt/birdstation/.venv/bin/birdstation-export
StandardOutput=journal
StandardError=journal
